# symphony.service
[Unit]
Description=Symphony Orchestrator
After=network.target docker.service

[Service]
Type=simple
User=marko
WorkingDirectory=/home/marko/projects/symphony
ExecStart=/home/marko/projects/symphony/start.sh
Restart=always
RestartSec=5
Environment=HOME=/home/marko
EnvironmentFile=-/home/marko/projects/symphony/.env

[Install]
WantedBy=multi-user.target
