# symphony.service
[Unit]
Description=Symphony Orchestrator
After=network.target

[Service]
Type=simple
User=marko
WorkingDirectory=/home/marko/Projects/symphony
ExecStart=/home/marko/Projects/symphony/bin/symphony --port 4000 WORKFLOW.md ENRICHMENT.md --i-understand-that-this-will-be-running-without-the-usual-guardrails
Restart=on-failure
RestartSec=5
Environment=PATH=/home/marko/.local/share/mise/installs/erlang/28.4.1/bin:/home/marko/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/home/marko/Projects/symphony/.env

[Install]
WantedBy=multi-user.target
