# symphony.service
[Unit]
Description=Symphony Orchestrator
After=network.target docker.service

[Service]
Type=simple
User=marko
WorkingDirectory=/home/marko/Projects/symphony
ExecStart=/home/marko/Projects/symphony/start.sh
Restart=always
RestartSec=5
Environment=HOME=/home/marko
EnvironmentFile=-/home/marko/Projects/symphony/.env

[Install]
WantedBy=multi-user.target
